# backend/crontab
0 0 * * * root NODE_ENV=production /usr/local/bin/node /app/bin/sitemap.js >> /var/log/cron.log 2>&1
